# deploy/launchd/io.sloptools.mcp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>io.sloptools.mcp</string>
  <key>ProgramArguments</key>
  <array>
    <string>@@SLOPTOOLS_BIN@@</string>
    <string>server</string>
    <string>--project-dir</string>
    <string>@@PROJECT_DIR@@</string>
    <string>--data-dir</string>
    <string>@@DATA_DIR@@</string>
    <string>--mcp-host</string>
    <string>127.0.0.1</string>
    <string>--mcp-port</string>
    <string>9420</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/sloptools-mcp.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/sloptools-mcp.log</string>
</dict>
</plist>
