# deploy/systemd/user/sloptools.service
[Unit]
Description=Sloptools Core MCP Server
After=network.target

[Service]
Type=simple
ExecStart=@@SLOPTOOLS_BIN@@ server --project-dir @@PROJECT_DIR@@ --data-dir @@DATA_DIR@@ --mcp-host 127.0.0.1 --mcp-port 9420
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
